# files/borg/borg.timer
[Unit]
Description=Run daily

[Timer]
OnBootSec=15min
OnUnitActiveSec=1hour
 
[Install]
WantedBy=timers.target
